# headphones.conf
[program:headphones]
user=sysadmin
directory=/opt/headphones
command=python Headphones.py --datadir=/config
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
